# graph-main/deploy/systemd/article-dz.socket
[Unit]
Description=ArticleDz Gunicorn socket

[Socket]
ListenStream=/run/article-dz.sock

[Install]
WantedBy=sockets.target
